# hplayer2@.service
[Unit]
Description=HPlayer
           
[Service]
Type=simple
WorkingDirectory=/opt/HPlayer2
ExecStartPre=-hplayer2-kill
ExecStart=/usr/local/bin/hplayer2 %i
ExecStop=/usr/local/bin/hplayer2-kill
Restart=always
           
[Install]
WantedBy=default.target
